# futuresonic-installer-debian/target/deb/DEBIAN/preinst
#! /bin/sh

set -e

# Stop FutureSonic service.
if [ -e /etc/init.d/futuresonic ]; then
  invoke-rc.d futuresonic stop
fi

# Backup database.
if [ -e /var/futuresonic/db ]; then
  rm -rf /var/futuresonic/db.backup
  cp -R /var/futuresonic/db /var/futuresonic/db.backup
fi

